# crontab
* * * * * /scripts/pipeline.sh >> /var/log/etl/pipeline.log 2>&1
